# checador.service
[Unit]
Description=Checador Fingerprint Time Clock
After=network.target

[Service]
Type=simple
User=arduino
Group=arduino
WorkingDirectory=/opt/checador
Environment="PATH=/usr/local/nbis/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 -m checador.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=checador

[Install]
WantedBy=multi-user.target